# mitamae/cookbooks/xkeysnail/templates/etc/systemd/system/xkeysnail.service.erb
[Unit]
Description=xkeysnail

[Service]
KillMode=process
ExecStart=/root/.local/bin/xkeysnail --quiet /home/kirikiriyamama/.xkeysnail
Type=simple
Restart=always
RestartSec=3s
Environment=DISPLAY=<%= node[:display] %>

[Install]
WantedBy=graphical.target
